# docker/all-in-one/supervisord.conf
[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
childlogdir=/data/logs

[program:samba]
command=/usr/sbin/samba --foreground --no-process-group
autostart=true
autorestart=true
priority=10
stdout_logfile=/data/logs/samba/samba-stdout.log
stderr_logfile=/data/logs/samba/samba-stderr.log

[program:mongodb]
command=/usr/local/bin/mongod --dbpath /data/mongodb --bind_ip 127.0.0.1 --port 27017 --logpath /data/logs/mongodb/mongod.log --quiet
user=zcloud
autostart=true
autorestart=true
priority=20
stdout_logfile=/data/logs/mongodb/mongod-stdout.log
stderr_logfile=/data/logs/mongodb/mongod-stderr.log

[program:webapp]
command=/bin/bash -c "source /etc/samba-conductor.env && exec node /home/zcloud/app/main.js"
directory=/home/zcloud/app
user=zcloud
autostart=true
autorestart=true
priority=30
startretries=10
startsecs=5
environment=NODE_ENV="production"
stdout_logfile=/data/logs/app/webapp-stdout.log
stderr_logfile=/data/logs/app/webapp-stderr.log
